# pi-ansible-config/roles/dev/files/rc.local
#!/bin/sh -e

/sbin/iw dev wlan0 set power_save off

exit 0
